# supervisord-sshd.conf
[program:sshd]
command=/run_sshd.sh
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0